# supervisord.conf
[supervisord]
nodaemon=true

[program:grafana-server]
user = grafana
command = sh /grafana.sh
